# src/main/profile.d/fn_path.sh
#!/bin/sh

export FN_PATH=/usr/share/fn
